# supervisor.conf
[program:ngircd]
user=irc
command=/usr/sbin/ngircd -f /etc/ngircd/ngircd.conf
redirect_stderr=true
autorestart=true
